# board/redpitaya/overlay/etc/profile.d/redpitaya_prompt.sh
export PS1="\[\e[0:31m\]red\[\e[m\]\[\e[1;31m\]pitaya\[\e[m\]\[\e[0;31m\]>\[\e[m\] "
